# contrib/harbor.service
[Unit]
Description=Harbor service registry
After=network.target caddy.service

[Service]
User=harbor
Environment=HARBOR_BACKEND=caddy
Environment=HARBOR_BACKEND_URL=unix:///run/caddy/admin.socket
Environment=HARBOR_STATIC_DIR=/etc/harbor/routes.d
ExecStart=gunicorn -k gthread --workers 1 --threads 16 harbor:app
Restart=on-failure

[Install]
WantedBy=multi-user.target